# config/crontab
0 9 * * * /usr/local/bin/node /app/retrieve.js
